# certbot/files/default/certsync/certsync.timer.j2
[Unit]
Description=Nightly run certsync

[Timer]
OnCalendar=*-*-* 01:00:00
RandomizedDelaySec=21600
Persistent=true

[Install]
WantedBy=timers.target
