# certbot/files/default/certsync/sync_certs.timer.j2
[Unit]
Description=Nightly run sync_certs

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=21600
Persistent=true

[Install]
WantedBy=timers.target
